# build/package/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:traffic]
autostart=true
autorestart=true
directory=/app
command=./traffic
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
